# Portable/Cygwin/usr/share/gvfs/mounts/nfs.mount
[Mount]
Type=nfs
Exec=/usr/libexec/gvfsd-nfs
AutoMount=false
Scheme=nfs
